# kadalu_operator/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:operator]
command=python3 /kadalu/main.py

[program:metrics]
command=python3 /kadalu/exporter.py
